# .config/autostart/mpd_pause_on_screensaver.desktop
[Desktop Entry]
Type=Application
Name=MPD pause on screensaver
Exec=mpd_pause_on_screensaver.sh
