# .config/autostart/xss-lock.desktop
[Desktop Entry]
Type=Application
Name=Lock screen on hibernation
Exec=xss-lock -- xdg-screensaver lock
